# deploy/hetzner/systemd/torus-ai-trading-strategy-hunt-broad.service
[Unit]
Description=TorusAiTrading Broad Strategy Hunt
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=torus
WorkingDirectory=/opt/TorusAiTrading
ExecStart=/usr/bin/bash /opt/TorusAiTrading/runtime/run-strategy-hunt-broad.sh

[Install]
WantedBy=multi-user.target
